# systemd/system/refresh-pacman-files.service
[Unit]
Description=Refresh pacman files database

[Service]
ExecStart=/usr/bin/pacman -Fy
Type=oneshot

[Install]
WantedBy=multi-user.target
