# systemd/system/refresh-pacman-files.timer
[Unit]
Description=Refresh pacman files database

[Timer]
OnBootSec=5min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
