# systemd/user/refresh-aur-updates.timer
[Unit]
Description=Refresh out-of-date AUR packages

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
