# systemd/user/refresh-pacman-updates.timer
[Unit]
Description=Refresh pacman out-of-date packages

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target