# systemd/user/refresh-aur-updates.service
[Unit]
Description=Refresh out-of-date AUR packages

[Service]
ExecStart=/home/rob/bin/refresh-aur-updates
Type=oneshot
Restart=on-failure
RestartSec=10s
RestartMaxDelaySec=5m
RestartSteps=5

[Install]
WantedBy=default.target
