# systemd/user/spothole.service
[Unit]
Description=Kill Spotify clients hanging around

[Service]
Type=oneshot
ExecStart=/home/rob/bin/spothole